# deploy/systemd/ayram-anomaly.service
[Unit]
Description=ML-Ayram Anomaly Detection
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=ayram
Group=ayram
WorkingDirectory=/home/ayram/ml-ayram

EnvironmentFile=/home/ayram/ml-ayram/.env

ExecStart=/home/ayram/ml-ayram/venv/bin/python3 \
    -m src.monitoring.anomaly_detector

TimeoutStartSec=600

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ayram-anomaly

NoNewPrivileges=true
PrivateTmp=true
