# deploy/systemd/ayram-dashboard.service
[Unit]
Description=ML-Ayram Dashboard (FastAPI + Uvicorn)
After=network-online.target postgresql.service
Wants=network-online.target
StartLimitIntervalSec=300s
StartLimitBurst=10

[Service]
Type=simple
User=ayram
Group=ayram
WorkingDirectory=/home/ayram/ml-ayram

EnvironmentFile=/home/ayram/ml-ayram/.env

ExecStart=/home/ayram/ml-ayram/venv/bin/uvicorn \
    src.dashboard.app:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 1 \
    --log-level warning

# CRITICO: always = reinicia SIEMPRE (crash, kill, reboot)
Restart=always
RestartSec=10s

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ayram-dashboard

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
