# deploy/systemd/ayram-train.service
[Unit]
Description=ML-Ayram Model Training — XGBoost + LSTM con Optuna
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=ayram
Group=ayram
WorkingDirectory=/home/ayram/ml-ayram

EnvironmentFile=/home/ayram/ml-ayram/.env

# Entrenar todos los pares/timeframes con optimización Optuna (30 trials)
# Primero actualiza features y labels frescos, luego entrena
ExecStartPre=/home/ayram/ml-ayram/venv/bin/python3 -m src.data.features
ExecStartPre=/home/ayram/ml-ayram/venv/bin/python3 -m src.data.labels
ExecStart=/home/ayram/ml-ayram/venv/bin/python3 \
    -m src.train --optimize --trials 30

# Notificar resultado por Telegram
ExecStartPost=/home/ayram/ml-ayram/venv/bin/python3 -c \
    "from src.notifications.telegram import send_message; \
     send_message('🔄 <b>Reentrenamiento semanal completado</b>\n\nModelos actualizados. Revisa logs para métricas.')"

# Timeout generoso: XGBoost+LSTM+Optuna puede tardar 4-6 horas
TimeoutStartSec=21600

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ayram-train

NoNewPrivileges=true
PrivateTmp=true
