# deploy/systemd/ayram-features.timer
[Unit]
Description=ML-Ayram Features + Labels — ejecutar cada hora
Requires=ayram-features.service

[Timer]
# Ejecutar cada 3 horas (minuto 10 para que el collector ya haya corrido)
# El job tarda ~2h con M15, así no se solapan ejecuciones
OnCalendar=*-*-* 00/3:10:00
AccuracySec=60s
Persistent=true

[Install]
WantedBy=timers.target
